# example/upsmonitor.service
[Unit]
Description=Upsmonitor Metrics Exporter
[Service]
Restart=on-failure
ExecStart=/root/upsmonitor/upsmonitor-prometheus
User=root
Group=root